# com.vectra.dev.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vectra.dev</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/franklintchakounteu/.gemini/antigravity/scratch/Vectra/scripts/start-dev.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/franklintchakounteu/.gemini/antigravity/scratch/Vectra/launchd.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/franklintchakounteu/.gemini/antigravity/scratch/Vectra/launchd.stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/franklintchakounteu/.gemini/antigravity/scratch/Vectra</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
